# iso/arch/i686/squashfs-root/var/lib/pacman/local/nbd-3.13-2/files
%FILES%
etc/
etc/nbd-server/
etc/nbd-server/config
usr/
usr/bin/
usr/bin/nbd-client
usr/bin/nbd-server
usr/bin/nbd-trdump
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nbd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/nbd.conf
usr/share/
usr/share/doc/
usr/share/doc/nbd/
usr/share/doc/nbd/README
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/nbd-server.1.gz
usr/share/man/man1/nbd-trdump.1.gz
usr/share/man/man5/
usr/share/man/man5/nbd-server.5.gz
usr/share/man/man5/nbdtab.5.gz
usr/share/man/man8/
usr/share/man/man8/nbd-client.8.gz

%BACKUP%
etc/nbd-server/config	a5b890148b80c5951dbe9daa2f487d63

