# iso/arch/i686/squashfs-root/var/lib/pacman/local/ndisc6-1.0.3-1/files
%FILES%
etc/
etc/rdnssd/
etc/rdnssd/merge-hook
etc/rdnssd/resolvconf-hook
usr/
usr/bin/
usr/bin/addr2name
usr/bin/dnssort
usr/bin/name2addr
usr/bin/ndisc6
usr/bin/rdisc6
usr/bin/rdnssd
usr/bin/rltraceroute6
usr/bin/tcpspray
usr/bin/tcpspray6
usr/bin/tcptraceroute6
usr/bin/tracert6
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rdnssd@.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/ndisc6.conf
usr/share/
usr/share/locale/
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/ndisc6.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/ndisc6.mo
usr/share/locale/en/
usr/share/locale/en/LC_MESSAGES/
usr/share/locale/en/LC_MESSAGES/ndisc6.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/ndisc6.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/ndisc6.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/addr2name.1.gz
usr/share/man/man1/dnssort.1.gz
usr/share/man/man1/name2addr.1.gz
usr/share/man/man1/tcpspray.1.gz
usr/share/man/man1/tcpspray6.1.gz
usr/share/man/man8/
usr/share/man/man8/ndisc6.8.gz
usr/share/man/man8/rdisc6.8.gz
usr/share/man/man8/rdnssd.8.gz
usr/share/man/man8/rltraceroute6.8.gz
usr/share/man/man8/tcptraceroute6.8.gz
usr/share/man/man8/tracert6.8.gz

