# iso/arch/i686/squashfs-root/var/lib/pacman/local/vpnc-0.5.3.svn550-1/files
%FILES%
etc/
etc/vpnc/
etc/vpnc/default.conf
etc/vpnc/vpnc-script
usr/
usr/bin/
usr/bin/cisco-decrypt
usr/bin/pcf2vpnc
usr/bin/vpnc
usr/bin/vpnc-disconnect
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/vpnc@.service
usr/share/
usr/share/doc/
usr/share/doc/vpnc/
usr/share/doc/vpnc/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/cisco-decrypt.1.gz
usr/share/man/man1/pcf2vpnc.1.gz
usr/share/man/man8/
usr/share/man/man8/vpnc.8.gz

%BACKUP%
etc/vpnc/default.conf	a3f4e0cc682f437e310a1c86ae198e45

