# iso/arch/i686/squashfs-root/var/lib/pacman/local/netctl-1.11-1/files
%FILES%
etc/
etc/ifplugd/
etc/ifplugd/netctl.action
etc/netctl/
etc/netctl/examples/
etc/netctl/examples/bonding
etc/netctl/examples/bridge
etc/netctl/examples/ethernet-custom
etc/netctl/examples/ethernet-dhcp
etc/netctl/examples/ethernet-static
etc/netctl/examples/macvlan-dhcp
etc/netctl/examples/macvlan-static
etc/netctl/examples/mobile_ppp
etc/netctl/examples/openvswitch
etc/netctl/examples/pppoe
etc/netctl/examples/tunnel
etc/netctl/examples/tuntap
etc/netctl/examples/vlan-dhcp
etc/netctl/examples/vlan-static
etc/netctl/examples/wireless-open
etc/netctl/examples/wireless-wep
etc/netctl/examples/wireless-wpa
etc/netctl/examples/wireless-wpa-config
etc/netctl/examples/wireless-wpa-configsection
etc/netctl/examples/wireless-wpa-static
etc/netctl/hooks/
etc/netctl/interfaces/
usr/
usr/bin/
usr/bin/netctl
usr/bin/netctl-auto
usr/bin/wifi-menu
usr/lib/
usr/lib/network/
usr/lib/network/auto.action
usr/lib/network/connections/
usr/lib/network/connections/README
usr/lib/network/connections/bond
usr/lib/network/connections/bridge
usr/lib/network/connections/dummy
usr/lib/network/connections/ethernet
usr/lib/network/connections/macvlan
usr/lib/network/connections/mobile_ppp
usr/lib/network/connections/openvswitch
usr/lib/network/connections/ppp
usr/lib/network/connections/pppoe
usr/lib/network/connections/tunnel
usr/lib/network/connections/tuntap
usr/lib/network/connections/vlan
usr/lib/network/connections/wireless
usr/lib/network/dhcp/
usr/lib/network/dhcp/README
usr/lib/network/dhcp/dhclient
usr/lib/network/dhcp/dhcpcd
usr/lib/network/globals
usr/lib/network/ip
usr/lib/network/network
usr/lib/network/rfkill
usr/lib/network/wpa
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/netctl-auto@.service
usr/lib/systemd/system/netctl-ifplugd@.service
usr/lib/systemd/system/netctl-sleep.service
usr/lib/systemd/system/netctl.service
usr/lib/systemd/system/netctl@.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/netctl
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/netctl-auto.1.gz
usr/share/man/man1/netctl.1.gz
usr/share/man/man5/
usr/share/man/man5/netctl.profile.5.gz
usr/share/man/man7/
usr/share/man/man7/netctl.special.7.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_netctl

