# iso/arch/i686/squashfs-root/var/lib/pacman/local/rfkill-0.5-1/files
%FILES%
usr/
usr/bin/
usr/bin/rfkill
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rfkill-block@.service
usr/lib/systemd/system/rfkill-unblock@.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-rfkill.rules
usr/share/
usr/share/licenses/
usr/share/licenses/rfkill/
usr/share/licenses/rfkill/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/rfkill.8.gz

