# iso/arch/x86_64/squashfs-root/var/lib/pacman/local/mdadm-3.4-1/files
%FILES%
etc/
etc/mdadm.conf
usr/
usr/bin/
usr/bin/mdadm
usr/bin/mdassemble
usr/bin/mdmon
usr/lib/
usr/lib/initcpio/
usr/lib/initcpio/hooks/
usr/lib/initcpio/hooks/mdadm
usr/lib/initcpio/install/
usr/lib/initcpio/install/mdadm
usr/lib/initcpio/install/mdadm_udev
usr/lib/systemd/
usr/lib/systemd/system-shutdown/
usr/lib/systemd/system-shutdown/mdadm.shutdown
usr/lib/systemd/system/
usr/lib/systemd/system/mdadm-grow-continue@.service
usr/lib/systemd/system/mdadm-last-resort@.service
usr/lib/systemd/system/mdadm-last-resort@.timer
usr/lib/systemd/system/mdmon@.service
usr/lib/systemd/system/mdmonitor.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/63-md-raid-arrays.rules
usr/lib/udev/rules.d/64-md-raid-assembly.rules
usr/share/
usr/share/man/
usr/share/man/man4/
usr/share/man/man4/md.4.gz
usr/share/man/man5/
usr/share/man/man5/mdadm.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/mdadm.8.gz
usr/share/man/man8/mdmon.8.gz

%BACKUP%
etc/mdadm.conf	5a37c112aa07dccdde62f9fa5b888607

