# iso/arch/i686/squashfs-root/var/lib/pacman/local/openssh-7.2p2-2/files
%FILES%
etc/
etc/pam.d/
etc/pam.d/sshd
etc/ssh/
etc/ssh/moduli
etc/ssh/ssh_config
etc/ssh/sshd_config
usr/
usr/bin/
usr/bin/findssl.sh
usr/bin/scp
usr/bin/sftp
usr/bin/ssh
usr/bin/ssh-add
usr/bin/ssh-agent
usr/bin/ssh-copy-id
usr/bin/ssh-keygen
usr/bin/ssh-keyscan
usr/bin/sshd
usr/lib/
usr/lib/ssh/
usr/lib/ssh/sftp-server
usr/lib/ssh/ssh-keysign
usr/lib/ssh/ssh-pkcs11-helper
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sshd.service
usr/lib/systemd/system/sshd.socket
usr/lib/systemd/system/sshd@.service
usr/lib/systemd/system/sshdgenkeys.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/sshd.conf
usr/share/
usr/share/licenses/
usr/share/licenses/openssh/
usr/share/licenses/openssh/LICENCE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/scp.1.gz
usr/share/man/man1/sftp.1.gz
usr/share/man/man1/slogin.1.gz
usr/share/man/man1/ssh-add.1.gz
usr/share/man/man1/ssh-agent.1.gz
usr/share/man/man1/ssh-copy-id.1.gz
usr/share/man/man1/ssh-keygen.1.gz
usr/share/man/man1/ssh-keyscan.1.gz
usr/share/man/man1/ssh.1.gz
usr/share/man/man5/
usr/share/man/man5/moduli.5.gz
usr/share/man/man5/ssh_config.5.gz
usr/share/man/man5/sshd_config.5.gz
usr/share/man/man8/
usr/share/man/man8/sftp-server.8.gz
usr/share/man/man8/ssh-keysign.8.gz
usr/share/man/man8/ssh-pkcs11-helper.8.gz
usr/share/man/man8/sshd.8.gz
var/
var/empty/

%BACKUP%
etc/ssh/ssh_config	7dac04c3896f2bbbffa7770ae4162966
etc/ssh/sshd_config	987f5acd4107718726563ce738449492
etc/pam.d/sshd	9d71ea52cd15e9d7c172b7b71cb03bc0

