# iso/arch/i686/squashfs-root/var/lib/pacman/local/smartmontools-6.5-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/smartd
etc/smartd.conf
usr/
usr/bin/
usr/bin/smartctl
usr/bin/smartd
usr/bin/update-smart-drivedb
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/smartd.service
usr/share/
usr/share/doc/
usr/share/doc/smartmontools/
usr/share/doc/smartmontools/AUTHORS
usr/share/doc/smartmontools/COPYING
usr/share/doc/smartmontools/ChangeLog
usr/share/doc/smartmontools/ChangeLog-5.0-6.0
usr/share/doc/smartmontools/INSTALL
usr/share/doc/smartmontools/NEWS
usr/share/doc/smartmontools/README
usr/share/doc/smartmontools/TODO
usr/share/doc/smartmontools/examplescripts/
usr/share/doc/smartmontools/examplescripts/Example1
usr/share/doc/smartmontools/examplescripts/Example2
usr/share/doc/smartmontools/examplescripts/Example3
usr/share/doc/smartmontools/examplescripts/Example4
usr/share/doc/smartmontools/examplescripts/Example5
usr/share/doc/smartmontools/examplescripts/Example6
usr/share/doc/smartmontools/examplescripts/README
usr/share/doc/smartmontools/smartd.conf
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/smartd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/smartctl.8.gz
usr/share/man/man8/smartd.8.gz
usr/share/man/man8/update-smart-drivedb.8.gz
usr/share/smartmontools/
usr/share/smartmontools/drivedb.h
usr/share/smartmontools/smartd_warning.d/
usr/share/smartmontools/smartd_warning.sh

%BACKUP%
etc/smartd.conf	624d033f31e2b301ffa67e4c028d9102
etc/conf.d/smartd	28c1b80c41dac52d433961dcfb4f70e0

