# iso/arch/x86_64/squashfs-root/var/lib/pacman/local/haveged-1.9.1-2/files
%FILES%
usr/
usr/bin/
usr/bin/haveged
usr/include/
usr/include/haveged/
usr/include/haveged/havege.h
usr/lib/
usr/lib/libhavege.so
usr/lib/libhavege.so.1
usr/lib/libhavege.so.1.1.0
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/haveged.service
usr/share/
usr/share/man/
usr/share/man/man3/
usr/share/man/man3/libhavege.3.gz
usr/share/man/man8/
usr/share/man/man8/haveged.8.gz

