# iso/arch/i686/squashfs-root/var/lib/pacman/local/dhcpcd-6.11.0-2/files
%FILES%
etc/
etc/dhcpcd.conf
usr/
usr/bin/
usr/bin/dhcpcd
usr/lib/
usr/lib/dhcpcd/
usr/lib/dhcpcd/dev/
usr/lib/dhcpcd/dev/udev.so
usr/lib/dhcpcd/dhcpcd-hooks/
usr/lib/dhcpcd/dhcpcd-hooks/01-test
usr/lib/dhcpcd/dhcpcd-hooks/02-dump
usr/lib/dhcpcd/dhcpcd-hooks/20-resolv.conf
usr/lib/dhcpcd/dhcpcd-hooks/30-hostname
usr/lib/dhcpcd/dhcpcd-run-hooks
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dhcpcd.service
usr/lib/systemd/system/dhcpcd@.service
usr/share/
usr/share/dhcpcd/
usr/share/dhcpcd/hooks/
usr/share/dhcpcd/hooks/10-wpa_supplicant
usr/share/dhcpcd/hooks/15-timezone
usr/share/dhcpcd/hooks/29-lookup-hostname
usr/share/licenses/
usr/share/licenses/dhcpcd/
usr/share/licenses/dhcpcd/LICENSE
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/dhcpcd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/dhcpcd-run-hooks.8.gz
usr/share/man/man8/dhcpcd.8.gz
var/
var/lib/
var/lib/dhcpcd/

%BACKUP%
etc/dhcpcd.conf	a685a737acd285b23664fbf24842c764

