# iso/arch/i686/squashfs-root/var/lib/pacman/local/dnsmasq-2.76-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/dnsmasq.conf
etc/dnsmasq.conf
usr/
usr/bin/
usr/bin/dnsmasq
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dnsmasq.service
usr/share/
usr/share/dnsmasq/
usr/share/dnsmasq/trust-anchors.conf
usr/share/locale/
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/dnsmasq.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/dnsmasq.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/dnsmasq.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/dnsmasq.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/dnsmasq.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/dnsmasq.mo
usr/share/locale/no/
usr/share/locale/no/LC_MESSAGES/
usr/share/locale/no/LC_MESSAGES/dnsmasq.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/dnsmasq.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/dnsmasq.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/dnsmasq.mo
usr/share/man/
usr/share/man/es/
usr/share/man/es/man8/
usr/share/man/es/man8/dnsmasq.8.gz
usr/share/man/fr/
usr/share/man/fr/man8/
usr/share/man/fr/man8/dnsmasq.8.gz
usr/share/man/man8/
usr/share/man/man8/dnsmasq.8.gz

%BACKUP%
etc/dnsmasq.conf	4b7ad4d68bfc2006e3f7f3293ec070ea

