# iso/arch/x86_64/squashfs-root/usr/lib/systemd/system/nfs-config.service
[Unit]
Description=Preprocess NFS configuration
After=local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/scripts/nfs-utils_env.sh
