# iso/arch/i686/squashfs-root/var/lib/pacman/local/darkhttpd-1.12-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/mimetypes
usr/
usr/bin/
usr/bin/darkhttpd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/darkhttpd.service
usr/share/
usr/share/licenses/
usr/share/licenses/darkhttpd/
usr/share/licenses/darkhttpd/LICENSE

%BACKUP%
etc/conf.d/mimetypes	d41d8cd98f00b204e9800998ecf8427e

