# iso/arch/x86_64/squashfs-root/var/lib/pacman/local/lvm2-2.02.154-3/files
%FILES%
etc/
etc/lvm/
etc/lvm/archive/
etc/lvm/backup/
etc/lvm/lvm.conf
etc/lvm/lvmlocal.conf
etc/lvm/profile/
etc/lvm/profile/cache-mq.profile
etc/lvm/profile/cache-smq.profile
etc/lvm/profile/command_profile_template.profile
etc/lvm/profile/metadata_profile_template.profile
etc/lvm/profile/thin-generic.profile
etc/lvm/profile/thin-performance.profile
usr/
usr/bin/
usr/bin/fsadm
usr/bin/lvchange
usr/bin/lvconvert
usr/bin/lvcreate
usr/bin/lvdisplay
usr/bin/lvextend
usr/bin/lvm
usr/bin/lvmchange
usr/bin/lvmconf
usr/bin/lvmconfig
usr/bin/lvmdiskscan
usr/bin/lvmdump
usr/bin/lvmetad
usr/bin/lvmsadc
usr/bin/lvmsar
usr/bin/lvreduce
usr/bin/lvremove
usr/bin/lvrename
usr/bin/lvresize
usr/bin/lvs
usr/bin/lvscan
usr/bin/pvchange
usr/bin/pvck
usr/bin/pvcreate
usr/bin/pvdisplay
usr/bin/pvmove
usr/bin/pvremove
usr/bin/pvresize
usr/bin/pvs
usr/bin/pvscan
usr/bin/vgcfgbackup
usr/bin/vgcfgrestore
usr/bin/vgchange
usr/bin/vgck
usr/bin/vgconvert
usr/bin/vgcreate
usr/bin/vgdisplay
usr/bin/vgexport
usr/bin/vgextend
usr/bin/vgimport
usr/bin/vgimportclone
usr/bin/vgmerge
usr/bin/vgmknodes
usr/bin/vgreduce
usr/bin/vgremove
usr/bin/vgrename
usr/bin/vgs
usr/bin/vgscan
usr/bin/vgsplit
usr/include/
usr/include/lvm2app.h
usr/include/lvm2cmd.h
usr/lib/
usr/lib/device-mapper/
usr/lib/device-mapper/libdevmapper-event-lvm2mirror.so
usr/lib/device-mapper/libdevmapper-event-lvm2raid.so
usr/lib/device-mapper/libdevmapper-event-lvm2snapshot.so
usr/lib/device-mapper/libdevmapper-event-lvm2thin.so
usr/lib/initcpio/
usr/lib/initcpio/hooks/
usr/lib/initcpio/hooks/lvm2
usr/lib/initcpio/install/
usr/lib/initcpio/install/lvm2
usr/lib/initcpio/install/sd-lvm2
usr/lib/initcpio/udev/
usr/lib/initcpio/udev/69-dm-lvm-metad.rules
usr/lib/libdevmapper-event-lvm2.so
usr/lib/libdevmapper-event-lvm2.so.2.02
usr/lib/libdevmapper-event-lvm2mirror.so
usr/lib/libdevmapper-event-lvm2raid.so
usr/lib/libdevmapper-event-lvm2snapshot.so
usr/lib/libdevmapper-event-lvm2thin.so
usr/lib/liblvm2app.so
usr/lib/liblvm2app.so.2.2
usr/lib/liblvm2cmd.so
usr/lib/liblvm2cmd.so.2.02
usr/lib/pkgconfig/
usr/lib/pkgconfig/lvm2app.pc
usr/lib/systemd/
usr/lib/systemd/system-generators/
usr/lib/systemd/system-generators/lvm2-activation-generator
usr/lib/systemd/system/
usr/lib/systemd/system/blk-availability.service
usr/lib/systemd/system/lvm2-lvmetad.service
usr/lib/systemd/system/lvm2-lvmetad.socket
usr/lib/systemd/system/lvm2-monitor.service
usr/lib/systemd/system/lvm2-pvscan@.service
usr/lib/systemd/system/sysinit.target.wants/
usr/lib/systemd/system/sysinit.target.wants/lvm2-lvmetad.socket
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/11-dm-lvm.rules
usr/lib/udev/rules.d/69-dm-lvm-metad.rules
usr/share/
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/lvm.conf.5.gz
usr/share/man/man7/
usr/share/man/man7/lvmcache.7.gz
usr/share/man/man7/lvmsystemid.7.gz
usr/share/man/man7/lvmthin.7.gz
usr/share/man/man8/
usr/share/man/man8/fsadm.8.gz
usr/share/man/man8/lvchange.8.gz
usr/share/man/man8/lvconvert.8.gz
usr/share/man/man8/lvcreate.8.gz
usr/share/man/man8/lvdisplay.8.gz
usr/share/man/man8/lvextend.8.gz
usr/share/man/man8/lvm-config.8.gz
usr/share/man/man8/lvm-dumpconfig.8.gz
usr/share/man/man8/lvm-lvpoll.8.gz
usr/share/man/man8/lvm.8.gz
usr/share/man/man8/lvm2-activation-generator.8.gz
usr/share/man/man8/lvmchange.8.gz
usr/share/man/man8/lvmconf.8.gz
usr/share/man/man8/lvmconfig.8.gz
usr/share/man/man8/lvmdiskscan.8.gz
usr/share/man/man8/lvmdump.8.gz
usr/share/man/man8/lvmetad.8.gz
usr/share/man/man8/lvmsadc.8.gz
usr/share/man/man8/lvmsar.8.gz
usr/share/man/man8/lvreduce.8.gz
usr/share/man/man8/lvremove.8.gz
usr/share/man/man8/lvrename.8.gz
usr/share/man/man8/lvresize.8.gz
usr/share/man/man8/lvs.8.gz
usr/share/man/man8/lvscan.8.gz
usr/share/man/man8/pvchange.8.gz
usr/share/man/man8/pvck.8.gz
usr/share/man/man8/pvcreate.8.gz
usr/share/man/man8/pvdisplay.8.gz
usr/share/man/man8/pvmove.8.gz
usr/share/man/man8/pvremove.8.gz
usr/share/man/man8/pvresize.8.gz
usr/share/man/man8/pvs.8.gz
usr/share/man/man8/pvscan.8.gz
usr/share/man/man8/vgcfgbackup.8.gz
usr/share/man/man8/vgcfgrestore.8.gz
usr/share/man/man8/vgchange.8.gz
usr/share/man/man8/vgck.8.gz
usr/share/man/man8/vgconvert.8.gz
usr/share/man/man8/vgcreate.8.gz
usr/share/man/man8/vgdisplay.8.gz
usr/share/man/man8/vgexport.8.gz
usr/share/man/man8/vgextend.8.gz
usr/share/man/man8/vgimport.8.gz
usr/share/man/man8/vgimportclone.8.gz
usr/share/man/man8/vgmerge.8.gz
usr/share/man/man8/vgmknodes.8.gz
usr/share/man/man8/vgreduce.8.gz
usr/share/man/man8/vgremove.8.gz
usr/share/man/man8/vgrename.8.gz
usr/share/man/man8/vgs.8.gz
usr/share/man/man8/vgscan.8.gz
usr/share/man/man8/vgsplit.8.gz

%BACKUP%
etc/lvm/lvm.conf	adcbadcb7e463efbcd761c2245cfef31

