# iso/arch/x86_64/squashfs-root/var/lib/pacman/local/nfs-utils-1.3.3-3/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/nfs-common.conf
etc/conf.d/nfs-server.conf
etc/exports
etc/exports.d/
etc/nfsmount.conf
etc/request-key.d/
etc/request-key.d/id_resolver.conf
usr/
usr/bin/
usr/bin/blkmapd
usr/bin/exportfs
usr/bin/mount.nfs
usr/bin/mount.nfs4
usr/bin/mountstats
usr/bin/nfsdcltrack
usr/bin/nfsidmap
usr/bin/nfsiostat
usr/bin/nfsstat
usr/bin/osd_login
usr/bin/rpc.gssd
usr/bin/rpc.idmapd
usr/bin/rpc.mountd
usr/bin/rpc.nfsd
usr/bin/rpc.statd
usr/bin/rpcdebug
usr/bin/showmount
usr/bin/sm-notify
usr/bin/start-statd
usr/bin/umount.nfs
usr/bin/umount.nfs4
usr/lib/
usr/lib/modules-load.d/
usr/lib/modules-load.d/nfs-utils.conf
usr/lib/systemd/
usr/lib/systemd/scripts/
usr/lib/systemd/scripts/nfs-utils_env.sh
usr/lib/systemd/system/
usr/lib/systemd/system/auth-rpcgss-module.service
usr/lib/systemd/system/nfs-blkmap.service
usr/lib/systemd/system/nfs-client.target
usr/lib/systemd/system/nfs-config.service
usr/lib/systemd/system/nfs-idmapd.service
usr/lib/systemd/system/nfs-mountd.service
usr/lib/systemd/system/nfs-server.service
usr/lib/systemd/system/nfs-utils.service
usr/lib/systemd/system/proc-fs-nfsd.mount
usr/lib/systemd/system/rpc-gssd.service
usr/lib/systemd/system/rpc-statd-notify.service
usr/lib/systemd/system/rpc-statd.service
usr/lib/systemd/system/rpc-svcgssd.service
usr/lib/systemd/system/var-lib-nfs-rpc_pipefs.mount
usr/share/
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/exports.5.gz
usr/share/man/man5/nfs.5.gz
usr/share/man/man5/nfsmount.conf.5.gz
usr/share/man/man7/
usr/share/man/man7/nfsd.7.gz
usr/share/man/man8/
usr/share/man/man8/blkmapd.8.gz
usr/share/man/man8/exportfs.8.gz
usr/share/man/man8/gssd.8.gz
usr/share/man/man8/idmapd.8.gz
usr/share/man/man8/mount.nfs.8.gz
usr/share/man/man8/mountd.8.gz
usr/share/man/man8/mountstats.8.gz
usr/share/man/man8/nfsd.8.gz
usr/share/man/man8/nfsdcltrack.8.gz
usr/share/man/man8/nfsidmap.8.gz
usr/share/man/man8/nfsiostat.8.gz
usr/share/man/man8/nfsstat.8.gz
usr/share/man/man8/rpc.gssd.8.gz
usr/share/man/man8/rpc.idmapd.8.gz
usr/share/man/man8/rpc.mountd.8.gz
usr/share/man/man8/rpc.nfsd.8.gz
usr/share/man/man8/rpc.sm-notify.8.gz
usr/share/man/man8/rpc.statd.8.gz
usr/share/man/man8/rpcdebug.8.gz
usr/share/man/man8/showmount.8.gz
usr/share/man/man8/sm-notify.8.gz
usr/share/man/man8/statd.8.gz
usr/share/man/man8/umount.nfs.8.gz
var/
var/lib/
var/lib/nfs/
var/lib/nfs/etab
var/lib/nfs/rmtab
var/lib/nfs/rpc_pipefs/
var/lib/nfs/sm.bak/
var/lib/nfs/sm/
var/lib/nfs/state
var/lib/nfs/v4recovery/
var/lib/nfs/xtab

%BACKUP%
etc/exports	e6ad3c7a59c7e4c24965a0e7da35026c
etc/nfsmount.conf	ce3b9fee8612087cbb69d46db34ce8ec
etc/conf.d/nfs-common.conf	95129639a19da1f262d4eb3ac73e0bdb
etc/conf.d/nfs-server.conf	7bd61e84ee05054f5232b2e8063ea08d

