# iso/arch/x86_64/squashfs-root/var/lib/pacman/local/openvpn-2.3.11-1/files
%FILES%
etc/
etc/openvpn/
usr/
usr/bin/
usr/bin/openvpn
usr/include/
usr/include/openvpn-plugin.h
usr/lib/
usr/lib/openvpn/
usr/lib/openvpn/plugins/
usr/lib/openvpn/plugins/openvpn-plugin-auth-pam.so
usr/lib/openvpn/plugins/openvpn-plugin-down-root.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/openvpn@.service
usr/share/
usr/share/doc/
usr/share/doc/openvpn/
usr/share/doc/openvpn/COPYING
usr/share/doc/openvpn/COPYRIGHT.GPL
usr/share/doc/openvpn/README
usr/share/doc/openvpn/README.IPv6
usr/share/doc/openvpn/README.auth-pam
usr/share/doc/openvpn/README.down-root
usr/share/doc/openvpn/README.polarssl
usr/share/doc/openvpn/management-notes.txt
usr/share/licenses/
usr/share/licenses/openvpn/
usr/share/licenses/openvpn/COPYING
usr/share/licenses/openvpn/COPYRIGHT.GPL
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/openvpn.8.gz
usr/share/openvpn/
usr/share/openvpn/contrib/
usr/share/openvpn/contrib/OCSP_check/
usr/share/openvpn/contrib/OCSP_check/OCSP_check.sh
usr/share/openvpn/contrib/README
usr/share/openvpn/contrib/multilevel-init.patch
usr/share/openvpn/contrib/openvpn-fwmarkroute-1.00/
usr/share/openvpn/contrib/openvpn-fwmarkroute-1.00/README
usr/share/openvpn/contrib/openvpn-fwmarkroute-1.00/fwmarkroute.down
usr/share/openvpn/contrib/openvpn-fwmarkroute-1.00/fwmarkroute.up
usr/share/openvpn/contrib/pull-resolv-conf/
usr/share/openvpn/contrib/pull-resolv-conf/client.down
usr/share/openvpn/contrib/pull-resolv-conf/client.up
usr/share/openvpn/examples/
usr/share/openvpn/examples/README
usr/share/openvpn/examples/client.conf
usr/share/openvpn/examples/firewall.sh
usr/share/openvpn/examples/home.up
usr/share/openvpn/examples/loopback-client
usr/share/openvpn/examples/loopback-server
usr/share/openvpn/examples/office.up
usr/share/openvpn/examples/openvpn-shutdown.sh
usr/share/openvpn/examples/openvpn-startup.sh
usr/share/openvpn/examples/server.conf
usr/share/openvpn/examples/static-home.conf
usr/share/openvpn/examples/static-office.conf
usr/share/openvpn/examples/tls-home.conf
usr/share/openvpn/examples/tls-office.conf
usr/share/openvpn/examples/xinetd-client-config
usr/share/openvpn/examples/xinetd-server-config

