# iso/arch/x86_64/squashfs-root/var/lib/pacman/local/xl2tpd-1.3.7-1/files
%FILES%
etc/
etc/ppp/
etc/ppp/options.l2tpd
etc/xl2tpd/
etc/xl2tpd/xl2tpd.conf.client-example
etc/xl2tpd/xl2tpd.conf.server-example
usr/
usr/bin/
usr/bin/pfc
usr/bin/xl2tpd
usr/bin/xl2tpd-control
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/xl2tpd.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/xl2tpd.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/pfc.1.gz
usr/share/man/man5/
usr/share/man/man5/l2tp-secrets.5.gz
usr/share/man/man5/xl2tpd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/xl2tpd.8.gz

%BACKUP%
etc/ppp/options.l2tpd	ed57eb67335cdf246251285dad502cc7

