# iso/arch/x86_64/squashfs-root/var/lib/pacman/local/dmraid-1.0.0.rc16.3-10/files
%FILES%
usr/
usr/bin/
usr/bin/dmevent_tool
usr/bin/dmraid
usr/include/
usr/include/dmraid/
usr/include/dmraid/display.h
usr/include/dmraid/dmraid.h
usr/include/dmraid/dmreg.h
usr/include/dmraid/format.h
usr/include/dmraid/lib_context.h
usr/include/dmraid/list.h
usr/include/dmraid/locking.h
usr/include/dmraid/metadata.h
usr/include/dmraid/misc.h
usr/include/dmraid/reconfig.h
usr/lib/
usr/lib/device-mapper/
usr/lib/device-mapper/libdmraid-events-isw.so
usr/lib/initcpio/
usr/lib/initcpio/hooks/
usr/lib/initcpio/hooks/dmraid
usr/lib/initcpio/install/
usr/lib/initcpio/install/dmraid
usr/lib/libdmraid-events-isw.so
usr/lib/libdmraid.so
usr/lib/libdmraid.so.1
usr/lib/libdmraid.so.1.0.0
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dmraid.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/dmraid.conf
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/dmevent_tool.8.gz
usr/share/man/man8/dmraid.8.gz

