# iso/arch/x86_64/squashfs-root/var/lib/pacman/local/rpcbind-0.2.3-3/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/rpcbind
usr/
usr/bin/
usr/bin/rpcbind
usr/bin/rpcinfo
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rpcbind.service
usr/lib/systemd/system/rpcbind.socket
usr/share/
usr/share/licenses/
usr/share/licenses/rpcbind/
usr/share/licenses/rpcbind/COPYING
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/rpcbind.8.gz
usr/share/man/man8/rpcinfo.8.gz
var/
var/lib/
var/lib/rpcbind/

%BACKUP%
etc/conf.d/rpcbind	d4c74d1be9f98344af138a15ad3b6f8c

