# iso/arch/x86_64/squashfs-root/var/lib/pacman/local/gssproxy-0.5.0-3/files
%FILES%
etc/
etc/gss/
etc/gss/mech.d/
etc/gss/mech.d/gssproxy.conf
etc/gssproxy/
etc/gssproxy/24-nfs-server.conf
etc/gssproxy/80-httpd.conf
etc/gssproxy/99-nfs-client.conf
etc/gssproxy/gssproxy.conf
usr/
usr/bin/
usr/bin/gssproxy
usr/lib/
usr/lib/gssproxy/
usr/lib/gssproxy/proxymech.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/gssproxy.service
usr/share/
usr/share/licenses/
usr/share/licenses/gssproxy/
usr/share/licenses/gssproxy/COPYING
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/gssproxy.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/gssproxy-mech.8.gz
usr/share/man/man8/gssproxy.8.gz
var/
var/lib/
var/lib/gssproxy/
var/lib/gssproxy/clients/
var/log/
var/log/gssproxy/

%BACKUP%
etc/gssproxy/gssproxy.conf	87302bfde98edf52e8c92f0002a575ce
etc/gssproxy/24-nfs-server.conf	215696bfcc2fe3ced2f250945cb45182
etc/gssproxy/80-httpd.conf	ebc5f2e55029fa77fb9732047f939d9a
etc/gssproxy/99-nfs-client.conf	b4cc264d7f04b480ed85b7dbc34c214a
etc/gss/mech.d/gssproxy.conf	f6e8373745bf5f7229ee878c0bf7d291

