# iso/arch/i686/squashfs-root/var/lib/pacman/local/iptables-1.6.0-1/files
%FILES%
etc/
etc/ethertypes
etc/iptables/
etc/iptables/empty.rules
etc/iptables/simple_firewall.rules
usr/
usr/bin/
usr/bin/arptables-compat
usr/bin/ebtables-compat
usr/bin/ip6tables
usr/bin/ip6tables-compat
usr/bin/ip6tables-compat-restore
usr/bin/ip6tables-compat-save
usr/bin/ip6tables-restore
usr/bin/ip6tables-save
usr/bin/iptables
usr/bin/iptables-compat
usr/bin/iptables-compat-restore
usr/bin/iptables-compat-save
usr/bin/iptables-restore
usr/bin/iptables-save
usr/bin/iptables-xml
usr/bin/xtables-compat-multi
usr/bin/xtables-multi
usr/include/
usr/include/libiptc/
usr/include/libiptc/ipt_kernel_headers.h
usr/include/libiptc/libip6tc.h
usr/include/libiptc/libiptc.h
usr/include/libiptc/libxtc.h
usr/include/libiptc/xtcshared.h
usr/include/xtables-version.h
usr/include/xtables.h
usr/lib/
usr/lib/iptables/
usr/lib/iptables/libarpt_mangle.so
usr/lib/iptables/libebt_802_3.so
usr/lib/iptables/libebt_ip.so
usr/lib/iptables/libebt_limit.so
usr/lib/iptables/libebt_log.so
usr/lib/iptables/libebt_mark.so
usr/lib/iptables/libebt_mark_m.so
usr/lib/iptables/libebt_nflog.so
usr/lib/iptables/libip6t_DNAT.so
usr/lib/iptables/libip6t_DNPT.so
usr/lib/iptables/libip6t_HL.so
usr/lib/iptables/libip6t_LOG.so
usr/lib/iptables/libip6t_MASQUERADE.so
usr/lib/iptables/libip6t_NETMAP.so
usr/lib/iptables/libip6t_REDIRECT.so
usr/lib/iptables/libip6t_REJECT.so
usr/lib/iptables/libip6t_SNAT.so
usr/lib/iptables/libip6t_SNPT.so
usr/lib/iptables/libip6t_ah.so
usr/lib/iptables/libip6t_dst.so
usr/lib/iptables/libip6t_eui64.so
usr/lib/iptables/libip6t_frag.so
usr/lib/iptables/libip6t_hbh.so
usr/lib/iptables/libip6t_hl.so
usr/lib/iptables/libip6t_icmp6.so
usr/lib/iptables/libip6t_ipv6header.so
usr/lib/iptables/libip6t_mh.so
usr/lib/iptables/libip6t_rt.so
usr/lib/iptables/libipt_CLUSTERIP.so
usr/lib/iptables/libipt_DNAT.so
usr/lib/iptables/libipt_ECN.so
usr/lib/iptables/libipt_LOG.so
usr/lib/iptables/libipt_MASQUERADE.so
usr/lib/iptables/libipt_NETMAP.so
usr/lib/iptables/libipt_REDIRECT.so
usr/lib/iptables/libipt_REJECT.so
usr/lib/iptables/libipt_SNAT.so
usr/lib/iptables/libipt_TTL.so
usr/lib/iptables/libipt_ULOG.so
usr/lib/iptables/libipt_ah.so
usr/lib/iptables/libipt_icmp.so
usr/lib/iptables/libipt_realm.so
usr/lib/iptables/libipt_ttl.so
usr/lib/iptables/libxt_AUDIT.so
usr/lib/iptables/libxt_CHECKSUM.so
usr/lib/iptables/libxt_CLASSIFY.so
usr/lib/iptables/libxt_CONNMARK.so
usr/lib/iptables/libxt_CONNSECMARK.so
usr/lib/iptables/libxt_CT.so
usr/lib/iptables/libxt_DSCP.so
usr/lib/iptables/libxt_HMARK.so
usr/lib/iptables/libxt_IDLETIMER.so
usr/lib/iptables/libxt_LED.so
usr/lib/iptables/libxt_MARK.so
usr/lib/iptables/libxt_NFLOG.so
usr/lib/iptables/libxt_NFQUEUE.so
usr/lib/iptables/libxt_NOTRACK.so
usr/lib/iptables/libxt_RATEEST.so
usr/lib/iptables/libxt_SECMARK.so
usr/lib/iptables/libxt_SET.so
usr/lib/iptables/libxt_SYNPROXY.so
usr/lib/iptables/libxt_TCPMSS.so
usr/lib/iptables/libxt_TCPOPTSTRIP.so
usr/lib/iptables/libxt_TEE.so
usr/lib/iptables/libxt_TOS.so
usr/lib/iptables/libxt_TPROXY.so
usr/lib/iptables/libxt_TRACE.so
usr/lib/iptables/libxt_addrtype.so
usr/lib/iptables/libxt_bpf.so
usr/lib/iptables/libxt_cgroup.so
usr/lib/iptables/libxt_cluster.so
usr/lib/iptables/libxt_comment.so
usr/lib/iptables/libxt_connbytes.so
usr/lib/iptables/libxt_connlimit.so
usr/lib/iptables/libxt_connmark.so
usr/lib/iptables/libxt_conntrack.so
usr/lib/iptables/libxt_cpu.so
usr/lib/iptables/libxt_dccp.so
usr/lib/iptables/libxt_devgroup.so
usr/lib/iptables/libxt_dscp.so
usr/lib/iptables/libxt_ecn.so
usr/lib/iptables/libxt_esp.so
usr/lib/iptables/libxt_hashlimit.so
usr/lib/iptables/libxt_helper.so
usr/lib/iptables/libxt_ipcomp.so
usr/lib/iptables/libxt_iprange.so
usr/lib/iptables/libxt_ipvs.so
usr/lib/iptables/libxt_length.so
usr/lib/iptables/libxt_limit.so
usr/lib/iptables/libxt_mac.so
usr/lib/iptables/libxt_mangle.so
usr/lib/iptables/libxt_mark.so
usr/lib/iptables/libxt_multiport.so
usr/lib/iptables/libxt_nfacct.so
usr/lib/iptables/libxt_osf.so
usr/lib/iptables/libxt_owner.so
usr/lib/iptables/libxt_physdev.so
usr/lib/iptables/libxt_pkttype.so
usr/lib/iptables/libxt_policy.so
usr/lib/iptables/libxt_quota.so
usr/lib/iptables/libxt_rateest.so
usr/lib/iptables/libxt_recent.so
usr/lib/iptables/libxt_rpfilter.so
usr/lib/iptables/libxt_sctp.so
usr/lib/iptables/libxt_set.so
usr/lib/iptables/libxt_socket.so
usr/lib/iptables/libxt_standard.so
usr/lib/iptables/libxt_state.so
usr/lib/iptables/libxt_statistic.so
usr/lib/iptables/libxt_string.so
usr/lib/iptables/libxt_tcp.so
usr/lib/iptables/libxt_tcpmss.so
usr/lib/iptables/libxt_time.so
usr/lib/iptables/libxt_tos.so
usr/lib/iptables/libxt_u32.so
usr/lib/iptables/libxt_udp.so
usr/lib/libip4tc.so
usr/lib/libip4tc.so.0
usr/lib/libip4tc.so.0.1.0
usr/lib/libip6tc.so
usr/lib/libip6tc.so.0
usr/lib/libip6tc.so.0.1.0
usr/lib/libiptc.so
usr/lib/libiptc.so.0
usr/lib/libiptc.so.0.0.0
usr/lib/libxtables.so
usr/lib/libxtables.so.11
usr/lib/libxtables.so.11.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/libip4tc.pc
usr/lib/pkgconfig/libip6tc.pc
usr/lib/pkgconfig/libiptc.pc
usr/lib/pkgconfig/xtables.pc
usr/lib/systemd/
usr/lib/systemd/scripts/
usr/lib/systemd/scripts/iptables-flush
usr/lib/systemd/system/
usr/lib/systemd/system/ip6tables.service
usr/lib/systemd/system/iptables.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/iptables-xml.1.gz
usr/share/man/man8/
usr/share/man/man8/ip6tables-restore.8.gz
usr/share/man/man8/ip6tables-save.8.gz
usr/share/man/man8/ip6tables.8.gz
usr/share/man/man8/iptables-extensions.8.gz
usr/share/man/man8/iptables-restore.8.gz
usr/share/man/man8/iptables-save.8.gz
usr/share/man/man8/iptables.8.gz
var/
var/lib/
var/lib/ip6tables/
var/lib/ip6tables/empty-filter.rules
var/lib/ip6tables/empty-mangle.rules
var/lib/ip6tables/empty-nat.rules
var/lib/ip6tables/empty-raw.rules
var/lib/ip6tables/empty-security.rules
var/lib/iptables/
var/lib/iptables/empty-filter.rules
var/lib/iptables/empty-mangle.rules
var/lib/iptables/empty-nat.rules
var/lib/iptables/empty-raw.rules
var/lib/iptables/empty-security.rules

