# iso/arch/i686/squashfs-root/var/lib/pacman/local/wpa_supplicant-1:2.5-3/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/wpa_supplicant.conf
etc/wpa_supplicant/
etc/wpa_supplicant/wpa_supplicant.conf
usr/
usr/bin/
usr/bin/wpa_cli
usr/bin/wpa_passphrase
usr/bin/wpa_supplicant
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/wpa_supplicant-nl80211@.service
usr/lib/systemd/system/wpa_supplicant-wired@.service
usr/lib/systemd/system/wpa_supplicant.service
usr/lib/systemd/system/wpa_supplicant@.service
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/fi.epitest.hostap.WPASupplicant.service
usr/share/dbus-1/system-services/fi.w1.wpa_supplicant1.service
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/wpa_supplicant.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/eapol_test.8.gz
usr/share/man/man8/wpa_background.8.gz
usr/share/man/man8/wpa_cli.8.gz
usr/share/man/man8/wpa_passphrase.8.gz
usr/share/man/man8/wpa_supplicant.8.gz

%BACKUP%
etc/wpa_supplicant/wpa_supplicant.conf	7fc85cd1fbfa4ce3e0c94ac63d1639ba

