# iso/arch/i686/squashfs-root/var/lib/pacman/local/partimage-0.6.9-9/files
%FILES%
etc/
etc/pam.d/
etc/pam.d/partimaged
etc/partimaged/
etc/partimaged/partimagedusers
usr/
usr/bin/
usr/bin/partimage
usr/bin/partimaged
usr/bin/partimaged-gencrt
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/partimaged.service
usr/share/
usr/share/doc/
usr/share/doc/partimage/
usr/share/doc/partimage/AUTHORS
usr/share/doc/partimage/BUGS
usr/share/doc/partimage/COPYING
usr/share/doc/partimage/ChangeLog
usr/share/doc/partimage/README
usr/share/doc/partimage/README.partimaged
usr/share/doc/partimage/partimage.lsm
usr/share/locale/
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/partimage.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/partimage.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/partimage.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/partimage.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/partimage.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/partimage.mo
usr/share/locale/gl/
usr/share/locale/gl/LC_MESSAGES/
usr/share/locale/gl/LC_MESSAGES/partimage.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/partimage.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/partimage.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/partimage.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/partimage.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/partimage.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/partimage.mo

%BACKUP%
etc/partimaged/partimagedusers	4181e75ffd76f9765607dc992e953e5a

